# deployment/config/s6/pimcore-automations/up
/etc/s6-overlay/scripts/pimcore-automations